# systemd.timer
[Unit]
Description=glassdoor jobs daily

[Timer]
OnCalendar=*-*-* 12:20:00
Persistent=true

[Install]
WantedBy=timers.target
